# public/config/captcha.conf
[program:captcha]
directory=/var/www/captcha
command=/var/www/captcha/captcha
autostart=true
autorestart=true
stderr_logfile=/var/log/captcha.err.log
stdout_logfile=/var/log/captcha.out.log
